# zpool_exporter/usr/local/etc/cron.d/zpool_exporter
*/1	*	*	*	*	root	/usr/local/sbin/zpool_exporter
