# smartmon_exporter/usr/local/etc/cron.d/smartmon_exporter
*/1	*	*	*	*	root	/usr/local/sbin/smartmon_exporter
