# lastlogin_exporter/usr/local/etc/cron.d/lastlogin_exporter
*/1	*	*	*	*	root	/usr/local/sbin/lastlogin_exporter /var/tmp/node_exporter/lastlogin_exporter.prom
